# etc/systemd/system/nix-daemon.service.d/override.conf
[Service]
Environment=TMPDIR=/var/tmp/nix-daemon
